# deploy/debian/moneyminder-guardian.service
[Unit]
Description=MoneyMinder health check and safe GitHub updater
Wants=network-online.target
After=network-online.target postgresql.service nginx.service tailscaled.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/moneyminder-guardian
ExecStart=/usr/local/sbin/moneyminder-guardian
TimeoutStartSec=240

[Install]
WantedBy=multi-user.target
